# packages/buendia-server/data/etc/init.d/buendia-warmup
#!/bin/sh
### BEGIN INIT INFO
# Provides: buendia-warmup
# Required-Start: $all
# Required-Stop:
# Should-Start:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Starts warming up the Buendia server once on boot.
# Description: Makes HTTP requests to the Buendia server until it comes up.
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin

case "$1" in
    start|restart|'')
        # This can take a long time and should not block system startup or
        # package installation, so it must run in the background.
        (/usr/bin/buendia-warmup -r localhost:9000 || service tomcat7 restart) 2>&1 | tee -a /var/log/buendia-warmup.log &
        ;;

    stop)
        # Nothing to do on shutdown.
        ;;

    *)
        echo "Usage: buendia-warmup [start|restart|stop]" >&2
        exit 3
        ;;
esac
