# packages/buendia-server/data/etc/cron.d/buendia-warmup
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root /usr/bin/buendia-warmup localhost:9000 >> /var/log/buendia-warmup.log 2>&1
