# packages/buendia-server/data/etc/cron.d/buendia-server
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

* * * * * root buendia-log buendia-warmup localhost:9000
