# packages/buendia-monitoring/data/etc/cron.d/buendia-monitoring
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

* * * * * root chmod -R a+rX /var/log
0 * * * * root buendia-log 'logrotate /etc/logrotate.conf; buendia-apply-limits'
10,20,30,40,50 * * * * root buendia-log buendia-apply-limits

0 */6 * * * root buendia-log 'du / | egrep -v "\s/sys|\s/proc"'
*/10 * * * * root buendia-log df
*/10 * * * * root buendia-log free
*/10 * * * * root buendia-log ps -efww
*/10 * * * * root buendia-log ip a
