# packages/buendia-backup/data/etc/cron.d/buendia-backup
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Back up to the first /dev/sd* partition we can find.
0 * * * * root . /usr/local/opt/buendia/utils.sh; if bool "$BACKUP_EXTERNAL"; then for dev in /dev/sd?[0-9]*; do [ -e $dev ] && buendia-log buendia-backup $dev $BACKUP_EXTERNAL_LIMIT_PERCENT && break; done; fi

# Back up to the internal directory.
10 * * * * root . /usr/local/opt/buendia/utils.sh; [ -n "$BACKUP_INTERNAL_DIR" ] && buendia-log buendia-backup $BACKUP_INTERNAL_DIR $BACKUP_INTERNAL_LIMIT_KB

# Restore from a USB drive
* * * * * root for dev in /dev/sd?[0-9]*; do [ -e $dev ] && buendia-log buendia-restore $dev && break; done;
