# packages/buendia-backup/data/etc/cron.d/buendia-backup
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 * * * * root /usr/bin/buendia-backup >> /var/log/buendia-backup.log 2>&1
