# packages/buendia-db-init/data/usr/share/buendia/systemd/buendia-db-init.service
[Unit]
Description=Buendia datbase reconfigure on boot
ConditionPathExists=!/etc/buendia-db-installed
Wants=mysql.service buendia-reconfigure.service
After=mysql.service
Before=buendia-reconfigure.service

[Service]
Type=oneshot
ExecStart=dpkg-reconfigure buendia-db-init

[Install]
WantedBy=multi-user.target
