# packages/buendia-db-init/data/usr/share/buendia/systemd/buendia-db-init.service
[Unit]
Description=Buendia datbase reconfigure on boot
ConditionPathExists=!/etc/buendia-db-installed
Wants=mysql.service
Wants=buendia-reconfigure.service
After=mysql.service
Before=buendia-reconfigure.service
Before=buendia-server.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/dpkg-reconfigure buendia-db-init

[Install]
WantedBy=multi-user.target
