# packages/buendia-security-monitor/data/lib/systemd/system/buendia-security-monitor.service
[Unit]
Description=Buendia security monitor
Wants=buendia-reconfigure.service
After=buendia-reconfigure.service

[Service]
Type=simple
ExecStart=/usr/bin/buendia-log buendia-security-monitor

[Install]
WantedBy=multi-user.target
