# platforms/edison/yocto/dnsmasq.service
[Unit]
Description=dnsmasq DNS and DHCP server

[Service]
Type=forking
PIDFile=/debian/run/dnsmasq/dnsmasq.pid
ExecStart=/usr/local/bin/enter-debian /etc/init.d/dnsmasq start
ExecStop=/usr/local/bin/enter-debian /etc/init.d/dnsmasq stop
Restart=on-failure
RestartSec=2
TimeoutSec=10

[Install]
WantedBy=multi-user.target
