# platforms/edison/yocto/mysql.service
[Unit]
Description=MySQL

[Service]
Type=forking
PIDFile=/debian/run/mysqld/mysqld.pid
ExecStart=/usr/local/bin/enter-debian /etc/init.d/mysql start
ExecStop=/usr/local/bin/enter-debian /etc/init.d/mysql stop
Restart=on-failure
RestartSec=2
TimeoutSec=10

[Install]
WantedBy=multi-user.target
