# packages/taskboard/data/etc/cron.d/taskboard
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * service taskboard status || service taskboard start
