# packages/buendia-networking/data/etc/cron.d/buendia-networking
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root buendia-update-hosts
* * * * * root sleep 30; buendia-log buendia-wifi-watchdog
