# platforms/edison-yocto/yocto/mount-sd.timer
[Unit]
Description=Periodically try to mount the SD card

[Timer]
OnBootSec=0
OnUnitActiveSec=120
Unit=mount-sd.service

[Install]
WantedBy=multi-user.target
