# platforms/edison-yocto/yocto/mysql.service
[Unit]
Description=MySQL
After=debian-mounts.service
Requires=debian-mounts.service

[Service]
Type=forking
PIDFile=/debian/run/mysqld/mysqld.pid
ExecStart=/usr/local/bin/enter-debian-nomount /etc/init.d/mysql start
ExecStop=/usr/local/bin/enter-debian-nomount /etc/init.d/mysql stop
Restart=on-failure
RestartSec=10
TimeoutSec=10

[Install]
WantedBy=multi-user.target
