# platforms/edison-yocto/yocto/openmrs-warmup.service
[Unit]
Description=OpenMRS warmup

# Tomcat and OpenMRS can take several minutes to start up; periodically making
# HTTP requests to them keeps them warm and responsive.

[Service]
Type=simple
ExecStart=/usr/bin/curl -v -u android:Android123 http://localhost:8080/openmrs/ws/rest/v1/projectbuendia/user
TimeoutStartSec=300
