# platforms/edison-yocto/yocto/openmrs-backup.timer
[Unit]
Description=Periodically back up the OpenMRS Database to an SD card


[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=openmrs-db-backup.service

[Install]
WantedBy=multi-user.target
