# platforms/edison-yocto/yocto/openmrs-warmup.timer
[Unit]
Description=Warms up OpenMRS periodically

[Timer]
OnBootSec=30
OnUnitActiveSec=300
Unit=openmrs-warmup.service

[Install]
WantedBy=multi-user.target
