# platforms/edison-yocto/yocto/tomcat7.service
[Unit]
Description=Apache Tomcat 7
After=debian-mounts.service
Requires=debian-mounts.service

[Service]
Type=forking
ExecStart=/usr/local/bin/enter-debian-nomount /etc/init.d/tomcat7 start
ExecStop=/usr/local/bin/enter-debian-nomount /etc/init.d/tomcat7 stop
Restart=on-failure
RestartSec=10
TimeoutSec=10

[Install]
WantedBy=multi-user.target
