# packages/buendia-tomcat7/data/etc/cron.d/buendia-tomcat7
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=

*/5 * * * * root buendia-log buendia-tomcat7-watchdog
