# packages/buendia-pkgserver/data/etc/cron.d/buendia-pkgserver
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root sleep 15; [ -e /dev/sda1 ] && buendia-log buendia-pkgserver-import /dev/sda1
* * * * * root sleep 45; [ -e /dev/sda1 ] && buendia-log buendia-pkgserver-import /dev/sda1
