# packages/buendia-server/control/postinst
#!/bin/bash

set -e; . /usr/share/buendia/utils.sh

case $1 in
    configure|abort-upgrade|abort-remove|abort-deconfigure)
        # Apply settings
        /usr/share/buendia/config.d/server

        # Restart Tomcat (necessary even if the settings haven't changed)
        service tomcat7 restart

        # Make the warmup script run once on startup
        update-rc.d buendia-warmup defaults
        ;;

    *) exit 1
esac
