# packages/buendia-server/control/postinst
#!/bin/bash

set -e

case $1 in
    configure|abort-upgrade|abort-remove|abort-deconfigure)
        # Apply settings
        /usr/share/buendia/config.d/server

        # Make the warmup script run once on startup
        update-rc.d buendia-warmup defaults
        ;;

    *) exit 1
esac
