# packages/buendia-dashboard/data/etc/cron.d/buendia-dashboard
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * * root /usr/bin/buendia-distill /var/log/large/requests /usr/share/buendia/distilled /var/log/requests

*/10 * * * * root /usr/bin/buendia-count /usr/share/buendia/counts

5,25,45 * * * * root cd /usr/share/buendia; zip -r dashboard/stats.zip counts distilled

# logrotate regularly runs once a day, but we need to run it more often on
# the /var/log/requests directory as the request logs can grow quickly.
0 * * * * root /usr/sbin/logrotate /etc/logrotate.d/buendia-dashboard
