# packages/buendia-update/control/postinst
#!/bin/bash

set -e; . /usr/share/buendia/utils.sh

case $1 in
    configure|abort-upgrade|abort-remove|abort-deconfigure)
        # Apply settings
        /usr/share/buendia/config.d/update

        # Make the update run once on startup
        update-rc.d buendia-autoupdate defaults

        # Ensure that cron is up so it can run the update periodically
        service cron start
        ;;

    *) exit 1
esac
