# packages/buendia-pkgclient/data/etc/cron.d/buendia-pkgclient
*/5 * * * * root /usr/bin/buendia-pkgclient-autoupdate -u 60
