# packages/buendia-pkgclient/data/etc/cron.d/buendia-pkgclient
*/5 * * * * root /etc/init.d/buendia-pkgclient-update start
