# packages/buendia-update/data/etc/cron.d/buendia-update
*/5 * * * * root /usr/bin/buendia-autoupdate -u 60
