# packages/buendia-update/data/etc/cron.d/buendia-update
*/5 * * * * root buendia-log buendia-autoupdate -u 60
