# packages/buendia-pkgclient/root/etc/cron.d/buendia-pkgclient.tab
*/5	*	*	*	*	root	/usr/bin/buendia-pkgclient-update > /dev/null
