# packages/buendia-pkgclient/root/etc/cron.d/buendia-pkgclient.tab
*/5	*	*	*	*	root	/usr/bin/pkgclient-update > /dev/null
