# packages/buendia-pkgclient/root/etc/cron.d/buendia-pkgclient.tab
*	*	*	*	*	root	/usr/bin/pkgclient-update > /dev/null
