# packages/buendia-utils/control/prerm
#!/bin/bash

set -e; . /usr/share/buendia/utils.sh

buendia-divert $1 /sbin/reboot
update-rc.d buendia-utils remove
