# packages/buendia-package-server/root/etc/profile.d/buendia_package_server.sh
#!/bin/bash
export DUSERVER_PACKAGE_DIR="/usr/share/buendia/packages"
export DUSERVER_PACKAGE_BASE_URL="packages.projectbuendia.org"
