# packages/buendia-utils/data/etc/cron.d/buendia-register-domains
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/1 * * * * root /usr/bin/buendia-register-domains >> /var/log/buendia-register-domains.log 2>&1
