# packages/buendia-update/data/etc/init.d/buendia-update
#!/bin/sh
### BEGIN INIT INFO
# Provides: buendia-autoupdate
# Required-Start: $all
# Required-Stop:
# Should-Start:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Keeps Buendia packages up to date.
# Description: Checks for and installs updates to all 'buendia-*' packages.
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin

case "$1" in
    start|restart|'')
        # This can take a long time and should not block system startup,
        # so it should run in the background.
        buendia-log buendia-autoupdate &
        ;;

    stop)
        # Nothing to do on shutdown.
        ;;

    *)
        echo "Usage: buendia-autoupdate [start|restart|stop]" >&2
        exit 3
        ;;
esac
