# packages/buendia-networking/control/postinst
#!/bin/bash

# This package installation takes place in the yocto environment. It enables
# and configures hostap and udhcp daemons. It sets the dns server to be located
# at the provided ip addres of this machine after the access point is up.
# Finally it signals apt that it requires a boot.
set -e

case $1 in
    configure)
        cat <<'EOF' | buendia-enter-yocto
hostapd_conf=/etc/hostapd/hostapd.conf
udhcpd_conf=/etc/hostapd/udhcpd-for-hostapd.conf

# Turn this machine into an access point.
systemctl enable hostapd
cp $hostapd_conf $hostapd_conf.pre-buendia
echo "ssid=buendia-demo" >> $hostapd_conf
echo "wpa_passphrase=buendia" >> $hostapd_conf
systemctl start hostapd

# Tell clients to use this machine's IP address for DNS.
ipaddr=''
while [ -z "$ipaddr" ]; do
  ipaddr=$(ifconfig | grep -A1 wlan0 | grep -o 'inet addr:[0-9][0-9.]*' | sed -e 's/.*://')
  sleep 1
done
cp $udhcpd_conf $udhcpd_conf.pre-buendia
echo "opt dns $ipaddr" >> $udhcpd_conf
EOF
        touch /var/run/reboot-required
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *) exit 1
esac

