# packages/buendia-networking/control/postinst
#!/bin/bash

# This package installation takes place in the yocto environment. It enables
# and configures hostap and udhcp daemons. It sets the dns server to be located
# at the provided ip addres of this machine after the access point is up.
# Finally it signals apt that it requires a boot.

set -e; . /usr/share/buendia/utils.sh

buendia-divert $1 /etc/dnsmasq.conf

case $1 in
    configure)
        cat <<'EOF' | buendia-enter-yocto
hostapd_conf=/etc/hostapd/hostapd.conf
udhcpd_conf=/etc/hostapd/udhcpd-for-hostapd.conf
udhcpd_service=/lib/systemd/system/udhcpd-for-hostapd.service
wpa_cli_actions=/etc/wpa_supplicant/wpa_cli-actions.sh

# Save a copy of the original configuration files.
divert() {
    [ -f "$1".pre-buendia ] || cp "$1" "$1".pre-buendia;
}
divert $hostapd_conf
divert $udhcpd_conf
divert $udhcpd_service
divert $wpa_cli_actions

# Don't use udhcpd with hostapd; we are using dnsmasq instead.
cat $udhcpd_service.pre-buendia | grep -v ExecStartPre | sed -e 's+ExecStart=.*+ExecStart=/bin/echo+' > $udhcpd_service
systemctl daemon-reload

# Make wpa_cli-actions.sh run only when wpa_supplicant_event is enabled.
cat <<WPA_CLI_ACTIONS_EOF > $wpa_cli_actions
#!/bin/sh
if systemctl -q is-enabled wpa_supplicant_event; then
    $wpa_cli_actions.pre-buendia "\$@"
else
    echo "event \$@ ignored because wpa_supplicant_event is disabled"
fi
WPA_CLI_ACTIONS_EOF
EOF
        # Apply the site settings.
        buendia-reconfigure networking

        # Also do network configuration on startup.
        update-rc.d buendia-networking defaults

        touch /var/run/reboot-required
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *) exit 1
esac

