# packages/buendia-networking/control/prerm
#!/bin/bash

set -e; . /usr/share/buendia/utils.sh

buendia-divert $1 /etc/dnsmasq.conf

case $1 in
    remove|upgrade|deconfigure)
        buendia-enter-yocto <<< '
            hostapd_conf=/etc/hostapd/hostapd.conf
            udhcpd_conf=/etc/hostapd/udhcpd-for-hostapd.conf
            udhcpd_service=/lib/systemd/system/udhcpd-for-hostapd.service

            # Restore the original configuration files in Yocto.
            mv $hostapd_conf.pre-buendia $hostapd_conf
            mv $udhcpd_conf.pre-buendia $udhcpd_conf
            mv $udhcpd_service.pre-buendia $udhcpd_service

            # Restore the original setup (DHCP client enabled, not an AP)
            systemctl disable hostapd
            systemctl stop hostapd || true
            systemctl enable wpa_supplicant_event
            systemctl start wpa_supplicant_event || true
        '
        touch /var/run/reboot-required
        ;;

    *) exit 1
esac
