# packages/buendia-pkgclient/control/postinst
#!/bin/bash
set -e

case $1 in
    configure|abort-upgrade|abort-remove|abort-deconfigure)
        /usr/share/buendia/config.d/pkgclient
        update-rc.d buendia-pkgclient-update defaults
        service cron restart
        ;;

    *) exit 1
esac
