# packages/buendia-utils/data/usr/share/buendia/systemd/buendia-reconfigure.service
[Unit]
Description=Buendia reconfigure on boot
ConditionPathExists=/etc/buendia-defer-reconfigure

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'rm -f /etc/buendia-defer-reconfigure; dpkg-reconfigure buendia-db-init; buendia-reconfigure'

[Install]
WantedBy=multi-user.target
