# packages/buendia-utils/data/usr/share/buendia/systemd/reboot-check.service
[Unit]
Description=Check for a reboot request

[Service]
Type=simple
ExecStart=/bin/sh -c 'if [ -e /tmp/reboot-requested ]; then rm -f /tmp/reboot-requested; /sbin/shutdown -r now "Remote Reboot Requested"; fi'

[Install]
WantedBy=multi-user.target
