# packages/buendia-utils/data/usr/share/buendia/systemd/reboot-check.timer
[Unit]
Description=Periodically check for a reboot request

[Timer]
OnBootSec=120
OnUnitActiveSec=2.5
AccuracySec=1
Unit=reboot-check.service

[Install]
WantedBy=multi-user.target
